# inittab
# /etc/inittab
::sysinit:/sbin/openrc sysinit
::wait:/sbin/openrc default
tty1::respawn:/usr/bin/checkra1n
tty2::respawn:/usr/sbin/usbmuxd -f
tty3::respawn:/bin/login -f root
tty4::respawn:/bin/login -f root
::ctrlaltdel:/sbin/reboot -f
